# bin/~/Documents/php/workspace/platform-api/service/kovey.socket
[Unit]
Description=Kovey Socket

[Socket]
ListenStream=9901
Accept=false
[Install]
WantedBy = sockets.target
