# bin/~/Documents/php/workspace/platform-service/service/kovey.service
[Unit]
Description=The Kovey Framwork Manager

[Service]
Type=forking
PIDFile=~/Documents/php/workspace/platform-service/run/kovey-rpc.pid
ExecStart=/usr/bin/php ~/Documents/php/workspace/platform-service/index.php
ExecStop=/bin/kill $MAINPID
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy = multi-user.target
